# conf/supervisor/conf.d/crontab.conf
[supervisord]
nodaemon=true

[program:system-cron]
command=cron -f
autostart=true
autorestart=true
stdout_logfile=/var/log/cron.log
