# conf/supervisor/conf.d/mysql.conf
[supervisord]
nodaemon=true

[program:system-mysql]
command=mysqld_safe
autostart=true
autorestart=true
stdout_logfile=/var/log/mysql.log
