# conf/supervisor/conf.d/phpfpm.conf
[supervisord]
nodaemon=true

[program:system-php7.2-fpm]
command=/usr/sbin/php-fpm7.2 -F
autostart=true
autorestart=true
priority=5
stdout_events_enabled=true
stderr_events_enabled=true
