# conf/supervisor/conf.d/phpfpm.conf
[supervisord]
nodaemon=true

[program:system-php7.4-fpm]
command=/usr/sbin/php-fpm7.4 -F
autostart=true
autorestart=true
stdout_logfile=/var/log/php-fpm7.4.log
