# conf/supervisor/conf.d/phpfpm.conf
[supervisord]
nodaemon=true

[program:php7.2-fpm]
command=php-fpm7.2 -F -c /etc/php/7.2/fpm/php-fpm.conf
autostart=true
autorestart=true
stdout_logfile=/var/log/php-fpm.log