# conf/crontab/crontab.conf
* * * * * ls -lah >> /tmp/crontab.log
* * * * * pwd >> /tmp/crontab-pwd.log