# conf/crontab/crontab.conf
## FIRST LINE NON DELETE!! ##
* * * * * php /var/www/artisan schedule:run >> /var/log/schedule-run.log
