# .docker/conf/crontab/crontab.conf
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

## CRON LISTS ##

* * * * ls -lah >> /tmp/crontab.log